# resources/serf.service
[Unit]
Description=cluster deployment system
After=network.target

[Service]
ExecStart=/usr/bin/serf agent -config-dir /etc/serf
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT
User=root
TimeoutStartSec=60
Restart=always

[Install]
WantedBy=multi-user.target
