# resources/serf.service
[Unit]
Description=cluster deployment system

[Service]
ExecStart=/usr/lib/redborder/bin/serf-agent.sh
KillSignal=SIGINT
User=root
TimeoutStartSec=60
Restart=on-failure

[Install]
WantedBy=multi-user.target
